# razorqt-resources/autostart/razor-qlipper-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=Qlipper
TryExec=qlipper
Exec=qlipper
OnlyShowIn=Razor;

#TRANSLATIONS_DIR=translations
