# razorqt-resources/autostart/razor-qlipper-autostart.desktop
[Desktop Entry]
Type=Application
Name=Qlipper
Name[sr]=Кјулипер
Name[sr@ijekavian]=Кјулипер
Name[sr@latin]=Kjuliper
Name[sr@ijekavianlatin]=Kjuliper
TryExec=qlipper
Exec=qlipper
OnlyShowIn=Razor;

Name[da]=Klippebord
Comment[da]=Qt klippebord
GenericName[da]=Klippebord
