# razorqt-resources/autostart/razor-qlipper-autostart.desktop
[Desktop Entry]
Type=Application
Name=Qlipper
TryExec=qlipper
Exec=qlipper
OnlyShowIn=Razor;

Name[da]=Klippebord
Comment[da]=Qt klippebord
GenericName[da]=Klippebord
