# razorqt-resources/autostart/razor-autosuspend.desktop.in
[Desktop Entry]
Type=Application
Version=1.0
Name=Razor Autosuspend
TryExec=razor-autosuspend
Exec=razor-autosuspend
OnlyShowIn=Razor;
X-Razor-Need-Tray=true
X-Razor-Module=true

#TRANSLATIONS_DIR=translations
