# razorqt-resources/autostart/razor-ptbatterysystemtray-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=PtBatterySystemTray
TryExec=ptbatterysystemtray
Exec=ptbatterysystemtray
OnlyShowIn=Razor;

#TRANSLATIONS_DIR=translations
