# razorqt-resources/autostart/razor-qlipper-autostart.desktop
[Desktop Entry]
Type=Application
Name=Qlipper
TryExec=qlipper
Exec=qlipper
OnlyShowIn=RAZOR;
