# razorqt-resources/autostart/razor-ptbatterysystemtray-autostart.desktop
[Desktop Entry]
Type=Application
Name=PtBatterySystemTray
TryExec=ptbatterysystemtray
Exec=ptbatterysystemtray
OnlyShowIn=RAZOR;
