# razorqt-resources/autostart/razor-qxkb-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=qxkb
TryExec=qxkb
Exec=qxkb
OnlyShowIn=Razor;

#TRANSLATIONS_DIR=translations
