# razorqt-resources/autostart/razor-qstardict-autostart.desktop
[Desktop Entry]
Type=Application
Name=QStarDict
TryExec=qstardict
Exec=qstardict
OnlyShowIn=RAZOR;
